# rocktalk.service
[Unit]
Description=RockTalk
After=network.target

[Service]
Type=simple
User=rocktalk
WorkingDirectory=/home/rocktalk/rocktalk-pi
ExecStart=/home/.nvm/versions/node/v24.15.0/bin/npm start
Restart=no
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target